# templates/server/restic-backup-notify@.service.j2
[Unit]
Description=Restic backup failure notifier (%i)
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/restic-backup-notify.py %i {{ restic_conf_directory }} server
