# templates/user/restic-backup-notify@.service.j2
[Unit]
Description=Restic backup failure notifier (%i)
After=network.target

[Service]
Type=oneshot
ExecStart=%h/.local/bin/restic-backup-notify.py %i %h/.config/restic-backup user
